# flybox/flybox.service
# systemd unit so the controller starts on boot.
# Install:
#   sudo cp flybox.service /etc/systemd/system/
#   sudo systemctl daemon-reload
#   sudo systemctl enable --now flybox
# Adjust User and WorkingDirectory to match your setup.

[Unit]
Description=FlyBox Controller (Pi 5 behavior box)
After=network-online.target
Wants=network-online.target

[Service]
User=pi5zms
WorkingDirectory=/home/pi5zms/flybox
# Use system python3 so picamera2 (apt-installed) is importable:
ExecStart=/usr/bin/python3 -m uvicorn app:app --host 0.0.0.0 --port 8000
Restart=on-failure
RestartSec=3

[Install]
WantedBy=multi-user.target
